# config/launchd/com.youtube.new-cohort-daily-channel-stats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube.new-cohort-daily-channel-stats</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.panels.daily_stats</string>
        <string>--channel-list</string>
        <string>data/channels/new_cohort/channel_ids.csv</string>
        <string>--mode</string>
        <string>channel</string>
        <string>--panel-name</string>
        <string>new_cohort</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/new_cohort_channel_stats_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/new_cohort_channel_stats_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
