# config/launchd/com.youtube.daily-discovery-non-intent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube.daily-discovery-non-intent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.collection.discover_non_intent</string>
        <string>--days-back</string>
        <string>1</string>
        <string>--strategies</string>
        <string>base,safesearch</string>
        <string>--skip-first-video</string>
        <string>--output</string>
        <string>data/channels/stream_a_prime/daily_discovery.csv</string>
        <string>--exclude-list</string>
        <string>data/channels/stream_a/daily_discovery.csv</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/daily_discovery_non_intent_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/daily_discovery_non_intent_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
