# config/launchd/com.youtube.gender-gap-enumeration.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube.gender-gap-enumeration</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.collection.enumerate_videos</string>
        <string>--channel-list</string>
        <string>data/channels/gender_gap/channel_ids.csv</string>
        <string>--output</string>
        <string>data/video_inventory/gender_gap_inventory.csv</string>
        <string>--max-runtime</string>
        <string>25200</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo</string>

    <key>StandardOutPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/gender_gap_enumeration_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/gender_gap_enumeration_stderr.log</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
