# scripts/com.youtube-longitudinal.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube-longitudinal.health-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.validation.health_check</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/health_check_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/katieapker/.youtube-longitudinal/repo/data/logs/health_check_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
